# seeder/config/supervisord.d/seeder.ini
[program:seeder]
command=/opt/bin/seed.sh
redirect_stderr=true
