# tracker/config/supervisord.d/tracker.ini
[program:tracker]
command=/opt/go/bin/chihaya -v=3 -config=/tmp/config.json -logtostderr=true
redirect_stderr=true
