# peer/config/supervisord.d/peer.ini
[program:peer]
command=/opt/bin/peer.sh
redirect_stderr=true
autorestart=false
